# SmartLocker_Pi5/smartlocker.service
[Unit]
Description=SmartLocker Pi Server
After=bluetooth.target network-online.target
Wants=network-online.target

[Service]
Type=simple
# UPDATE THESE VALUES:
Environment="TTGO_MAC=14:2B:2F:A7:94:22"
Environment="WEBHOOK_URL=https://script.google.com/macros/s/YOUR_DEPLOYMENT_ID/exec"

ExecStart=/opt/smartlocker/start_locker.sh ${TTGO_MAC} ${WEBHOOK_URL}
Restart=always
RestartSec=10

# Logging
StandardOutput=append:/var/log/smartlocker.log
StandardError=append:/var/log/smartlocker.log

[Install]
WantedBy=multi-user.target
